# tp3/systemd/units/Backup.timer
[Unit]
Description=Lance le service backup toute les heures

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target